# deploy/systemd/floodcast-minxiong-api.service
[Unit]
Description=FloodCastMinxiong observation API and operator view
Documentation=https://github.com/KageRyo/FloodCastMinxiong
After=network.target

[Service]
Type=simple
User=floodcast-minxiong
Group=floodcast-minxiong
WorkingDirectory=/opt/floodcast-minxiong
EnvironmentFile=-/etc/floodcast-minxiong.env
StateDirectory=floodcast-minxiong
ExecStart=/opt/floodcast-minxiong/.venv/bin/floodcast-minxiong-serve \
  --store /var/lib/floodcast-minxiong/operations \
  --host 127.0.0.1 \
  --port 8080
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/var/lib/floodcast-minxiong

[Install]
WantedBy=multi-user.target
