# deploy/systemd-user/floodcast-minxiong-alert-receiver.service
[Unit]
Description=FloodCastMinxiong Alertmanager webhook audit receiver
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/share/floodcast-minxiong/venv/bin/floodcast-minxiong-alert-receiver \
  --host 127.0.0.1 \
  --port 9087 \
  --output %h/.local/share/floodcast-minxiong/notifications/alerts.jsonl
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
