# deploy/systemd-user/floodcast-minxiong-alertmanager.service
[Unit]
Description=FloodCastMinxiong Alertmanager
After=network.target floodcast-minxiong-alert-receiver.service
Wants=floodcast-minxiong-alert-receiver.service

[Service]
Type=simple
ExecStart=%h/.local/share/floodcast-minxiong/bin/alertmanager \
  --config.file=%h/.local/share/floodcast-minxiong/config/prometheus/alertmanager.yml \
  --storage.path=%h/.local/share/floodcast-minxiong/alertmanager \
  --web.listen-address=127.0.0.1:9093
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
