# deploy/systemd-user/floodcast-minxiong-runner.service
[Unit]
Description=GitHub Actions runner for FloodCastMinxiong host-bound contracts
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/floodcast-minxiong/actions-runner
ExecStart=%h/.local/share/floodcast-minxiong/actions-runner/run.sh
Restart=always
RestartSec=10
KillMode=process
TimeoutStopSec=5min
NoNewPrivileges=true

[Install]
WantedBy=default.target
