# deploy/systemd-user/floodcast-minxiong-api.service
[Unit]
Description=FloodCastMinxiong observation API and operator view
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.config/floodcast-minxiong/env
ExecStart=%h/.local/share/floodcast-minxiong/venv/bin/floodcast-minxiong-serve \
  --store %h/.local/share/floodcast-minxiong/operations \
  --host 127.0.0.1 \
  --port 8180
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
