# deploy/systemd-user/floodcast-minxiong-shadow.timer
[Unit]
Description=Evaluate FloodCastMinxiong shadow gate hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120
Unit=floodcast-minxiong-shadow.service

[Install]
WantedBy=timers.target
