# deploy/systemd-user/floodcast-minxiong-collector.timer
[Unit]
Description=Collect FloodCastMinxiong observations every 10 minutes

[Timer]
OnCalendar=*-*-* *:0/10:00
Persistent=true
RandomizedDelaySec=30
Unit=floodcast-minxiong-collector.service

[Install]
WantedBy=timers.target
