# deploy/systemd-user/minxiong-hydrocast-collector.service
[Unit]
Description=MinxiongHydroCast live observation collection
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/minxiong-hydrocast/env
ExecStart=%h/.local/share/minxiong-hydrocast/venv/bin/minxiong-hydrocast-operations \
  --mode live \
  --once \
  --alert-source api \
  --rain-source api \
  --flood-source api \
  --flood-max-age-minutes 90 \
  --api-timeout-seconds 30 \
  --store %h/.local/share/minxiong-hydrocast/operations \
  --debug-dir %h/.local/share/minxiong-hydrocast/debug \
  --summary-output %h/.local/share/minxiong-hydrocast/run_summaries/operations.json \
  --log-output %h/.local/share/minxiong-hydrocast/run_logs.jsonl
NoNewPrivileges=true
PrivateTmp=true
