# deploy/systemd-user/minxiong-hydrocast-event-discover.service
[Unit]
Description=MinxiongHydroCast continuous radar-event evidence discovery
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/minxiong-hydrocast/env
WorkingDirectory=%h/.local/share/minxiong-hydrocast
ExecStart=%h/.local/share/minxiong-hydrocast/venv/bin/minxiong-hydrocast-event-discover \
  --repository-root %h/.local/share/minxiong-hydrocast \
  --initial-lookback-minutes 120 \
  --before-minutes 60 \
  --after-minutes 60 \
  --max-candidate-window-minutes 480 \
  --event-threshold-dbz 35 \
  --summary-output %h/.local/share/minxiong-hydrocast/run_summaries/event_discover.json \
  --log-output %h/.local/share/minxiong-hydrocast/run_logs.jsonl
NoNewPrivileges=true
PrivateTmp=true
